# scripts/pisugar/systemd/pisugar-watchdog.service
[Unit]
Description=PiSugar 3 MCU watchdog
After=pisugar-server.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pisugar-watchdog.sh
